# systemd-examples/feeds-to-readeck.timer
[Unit]
Description=Timer for Feeds to Readeck

[Timer]
OnStartupSec=1m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
